# deploy/pigpiod.service
[Unit]
Description=Daemon required to control GPIO pins via /dev/gpiomem
After=network.target

[Service]
Type=simple
# pigpiod installs to /usr/local/bin/pigpiod when compiled from source
# Note: pigpiod doesn't fork or create PID file, so Type=simple (not forking)
ExecStart=/usr/local/bin/pigpiod
ExecStop=/bin/killall pigpiod
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target

