# deploy/silvia-pid.service
[Unit]
Description=Silvia PID Controller with Docker Compose
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/silvia-pid
EnvironmentFile=-/opt/silvia-pid/.env

# Start command - start containers (health checks may take time, but containers will start)
ExecStart=/usr/bin/docker compose -f /opt/silvia-pid/docker-compose.yml up -d

# Stop command
ExecStop=/usr/bin/docker compose down

# Restart policy
Restart=on-failure
RestartSec=10s

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=silvia-pid

[Install]
WantedBy=multi-user.target

